# .service
[Unit]
Description=Shinshi Service
After=network.target

[Service]
User=user
WorkingDirectory=/usr/shinshi/
ExecStart=/usr/bin/uv run -- python3 -OO -m shinshi
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
